# deploy/systemd/clankermux.service.d/00-verify-deps.conf
# /etc/systemd/system/clankermux.service.d/00-verify-deps.conf
# Verify node_modules matches the committed, integrity-hashed bun.lock before
# the proxy (and the dashboard build) start. Named 00- so this ExecStartPre is
# ordered ahead of dashboard-build.conf. Fail-closed: a mismatch aborts startup.
[Service]
ExecStartPre=/home/darken/clankermux/scripts/verify-deps.sh
