# deploy/systemd/clankermux.service.d/backend-port.conf
# /etc/systemd/system/clankermux.service.d/backend-port.conf
# Move the app behind the Caddy front proxy (deploy/caddy/): the app binds
# 127.0.0.1:8090 while Caddy owns the client-facing :8080, holding new
# connections during app restarts so clients never see ECONNREFUSED. The
# client-facing port stays :8080 — only the process behind it changes.
# Binding to loopback also removes the management-API network exposure warning
# (the dashboard/API is no longer directly reachable from the LAN).
# ⚠️ ONLY apply this drop-in together with the Caddy front proxy — with the
# app on loopback:8090 and no Caddy on :8080, clients reach nothing.
[Service]
Environment=PORT=8090
# NOTE: the bind host is read via readEnv("HOST") in packages/core/src/env.ts,
# which only checks the prefixed forms (CLANKERMUX_HOST / BETTER_CCFLARE_HOST /
# ccflare_HOST) — a bare HOST= is silently ignored and the app binds 0.0.0.0.
Environment=CLANKERMUX_HOST=127.0.0.1
