# deploy/systemd/clankermux.service.d/dashboard-build.conf
# /etc/systemd/system/clankermux.service.d/dashboard-build.conf
# Regenerate the gitignored inline DB workers and rebuild the dashboard before
# the proxy starts, so both self-heal on every restart. Each step runs through a
# content-hash guard (scripts/guarded-build.ts): it hashes the source inputs AND
# the produced artifacts (content, not mtime) and skips the expensive build when
# nothing changed — cutting the restart window from ~30-90s to a few seconds.
# It falls back to a full unconditional build on first run, a missing/corrupt
# marker, or a missing/stale artifact (an inline worker still holding its empty
# EMBEDDED_* placeholder counts as stale). A failed build never writes a marker,
# so it re-runs next restart — the ExecStartPre stays fail-closed and blocking.
# db-workers runs first because the proxy imports the inline workers at startup.
[Service]
ExecStartPre=/home/darken/.bun/bin/bun run build:db-workers:guarded
ExecStartPre=/home/darken/.bun/bin/bun run build:dashboard:guarded
