# deploy/systemd/clankermux.service.d/hardening.conf
# /etc/systemd/system/clankermux.service.d/hardening.conf
# Sandbox hardening for the network-facing proxy. Compatible with a home-dir
# source install: ProtectHome is NOT set (the working tree + DB live under
# /home/darken), and MemoryDenyWriteExecute is NOT set (bun's JIT needs W+X).
[Service]
# Whole filesystem read-only except the few paths the service must write.
ProtectSystem=strict
ReadWritePaths=/home/darken/clankermux /home/darken/.config/clankermux /home/darken/.bun /home/darken/.cache
# Kernel / host surface the proxy never needs.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
# Networking: only the families a TCP proxy uses (AF_NETLINK is needed by
# glibc getaddrinfo for DNS / interface enumeration).
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Drop all capabilities — binds port 8080 (>1024), so none are required.
CapabilityBoundingSet=
AmbientCapabilities=
# Syscall allowlist for typical services. (Subtracting ~@privileged/@resources
# on top of this was verified to SIGSYS-kill bun's runtime, so it is omitted.)
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Default-deny world-readable files created by the service.
UMask=0077
