# deploy/systemd/clankermux.service.d/stop-timeout.conf
# /etc/systemd/system/clankermux.service.d/stop-timeout.conf
# Give the in-app shutdown watchdog (SHUTDOWN_WATCHDOG_MS = 300s in
# apps/server/src/server.ts) room to run: systemd's default TimeoutStopSec is
# 90s, which would SIGKILL the process mid-drain long before the watchdog
# fires. 330s = watchdog + margin, so the app always exits on its own terms
# (clean drain or watchdog exit) and systemd never has to escalate to SIGKILL.
# Safe because the Caddy front proxy (deploy/caddy/) holds new connections for
# the whole drain — a long stop is invisible to clients.
[Service]
TimeoutStopSec=330
